# app-emulation/hv-tools/files/hv_kvp_daemon.initd
#!/sbin/openrc-run

command=/usr/sbin/hv_kvp_daemon
command_args="-n"

name="hv_kvp_daemon"
description="Hyper-V KVP Protocol Daemon"
