# app-emulation/hv-tools/files/hv_vss_daemon.initd
#!/sbin/openrc-run

command=/usr/sbin/hv_vss_daemon
command_args="-n"

name="hv_vss_daemon"
description="Hyper-V VSS Protocol Daemon"
